# duplex-scan-for-any-printer/rootfs/etc/s6-overlay/s6-rc.d/scan-agent/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Scan Agent
# Config is at /data/config.yaml — generated from HA options by init-prepare
# Exposes internal API on 127.0.0.1:8098
# Auto-restarts if crashes
# ==============================================================================

if command -v bashio &>/dev/null; then
    bashio::log.info "Starting Scan Agent..."
else
    echo "[scan-agent] Starting Scan Agent..."
fi

cd /app

exec python3 -m src.main --config /data/config.yaml
