# duplex-scan-for-any-printer/deploy/scan-agent.service
[Unit]
Description=Scan Agent Service
Documentation=https://github.com/yourusername/PrinterDuplexScanForCheapPrint
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=scanagent
Group=scanagent

# Working directory
WorkingDirectory=/opt/scan-agent

# Environment
Environment="PATH=/opt/scan-agent/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="SCAN_INBOX_BASE=/scan_inbox"
Environment="SCAN_OUTPUT_DIR=/scan_out"
Environment="PYTHONPATH=/opt/scan-agent"

# Start command
ExecStart=/opt/scan-agent/.venv/bin/python -m src.main --config /opt/scan-agent/config.yaml

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
# Adjust these based on your hardware
LimitNOFILE=65536
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scan-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/scan_inbox /scan_out /var/log/scan-agent

# Required directories and permissions
RuntimeDirectory=scan-agent
RuntimeDirectoryMode=0755
LogsDirectory=scan-agent
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target
