# duplex-scan-for-any-printer/rootfs/etc/s6-overlay/s6-rc.d/web-ui/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Web UI stopped"
