# duplex-scan-for-any-printer/rootfs/etc/s6-overlay/s6-rc.d/web-ui/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Web UI Server
# Serves the dashboard on port 8099; proxies /api/* to agent API on 8098
# ==============================================================================

if command -v bashio &>/dev/null; then
    bashio::log.info "Starting Web UI on port 8099..."
else
    echo "[web-ui] Starting Web UI on port 8099..."
fi

cd /app

exec python3 -m src.web_ui_server
