# duplex-scan-for-any-printer/rootfs/etc/s6-overlay/s6-rc.d/init-prepare/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Prepare environment and configuration
# Runs once before services start
#
# Two modes:
#   HAOS  — bashio available: reads addon options, generates /data/config.yaml
#   Docker — bashio absent:   /data/config.yaml already mounted, used as-is
# ==============================================================================

set -e

# ── Logging helper ───────────────────────────────────────────────────────────
log_info() {
    if command -v bashio &>/dev/null; then
        bashio::log.info "$1"
    else
        echo "[init-prepare] $1"
    fi
}

log_info "Preparing Scan Agent environment..."

# ── Read configuration ───────────────────────────────────────────────────────
if command -v bashio &>/dev/null; then
    # HAOS mode: read supervisor options and generate /data/config.yaml
    SESSION_TIMEOUT=$(bashio::config 'session_timeout_seconds' '300')
    DELETE_INBOX=$(bashio::config 'delete_inbox_files_after_process' 'true')
    TEST_MODE=$(bashio::config 'test_mode' 'false')
    PRINTER_ENABLED=$(bashio::config 'printer.enabled' 'false')
    PRINTER_NAME=$(bashio::config 'printer.name' '')
    PRINTER_IP=$(bashio::config 'printer.ip' '')
    FTP_USERNAME=$(bashio::config 'ftp.username' '')
    FTP_PASSWORD=$(bashio::config 'ftp.password' '')
    TG_ENABLED=$(bashio::config 'telegram.enabled' 'false')
    TG_TOKEN=$(bashio::config 'telegram.bot_token' '')
    TG_AUTH=$(bashio::config 'telegram.authorized_users' '' | jq -Rc '[splits(",") | gsub("^ +| +$";"") | select(length > 0) | tonumber]' 2>/dev/null || echo '[]')
    TG_NOTIFY=$(bashio::config 'telegram.notify_chat_ids' '' | jq -Rc '[splits(",") | gsub("^ +| +$";"") | select(length > 0) | tonumber]' 2>/dev/null || echo '[]')
    TG_NOTIFY_ON_READY=$(bashio::config 'telegram.notify_on_session_ready' 'true')

    log_info "Log level: $(bashio::config 'log_level' 'info')"
    log_info "Generating /data/config.yaml..."

    cat > /data/config.yaml << EOF
# Auto-generated from Home Assistant addon options — DO NOT EDIT MANUALLY
inbox_base: /share/scan_inbox
subdirs:
  scan_duplex: scan_duplex
  copy_duplex: copy_duplex
  scan_document: scan_document
  card_2in1: card_2in1
  confirm: confirm
  confirm_print: confirm_print
  reject: reject
  test_print: test_print
output_dir: /share/scan_out
session_timeout_seconds: ${SESSION_TIMEOUT}
a4_page:
  width_pt: 595
  height_pt: 842
margin_pt: 15
delete_inbox_files_after_process: ${DELETE_INBOX}
test_mode: ${TEST_MODE}

printer:
  enabled: ${PRINTER_ENABLED}
  name: "${PRINTER_NAME}"
  ip: "${PRINTER_IP}"

telegram:
  enabled: ${TG_ENABLED}
  bot_token: "${TG_TOKEN}"
  authorized_users: ${TG_AUTH}
  notify_chat_ids: ${TG_NOTIFY}
  notify_on_session_ready: ${TG_NOTIFY_ON_READY}
EOF

    log_info "/data/config.yaml written."
    log_info "--- Generated config ---"
    cat /data/config.yaml
    log_info "--- End config ---"
    log_info "telegram.authorized_users (raw): $(bashio::config 'telegram.authorized_users' '')"
    log_info "telegram.authorized_users (parsed): ${TG_AUTH}"
    log_info "telegram.notify_chat_ids (raw): $(bashio::config 'telegram.notify_chat_ids' '')"
    log_info "telegram.notify_chat_ids (parsed): ${TG_NOTIFY}"
else
    # Plain Docker mode: /data/config.yaml must already be mounted from test_data/
    if [ ! -f "/data/config.yaml" ]; then
        echo "[init-prepare] ERROR: /data/config.yaml not found."
        echo "[init-prepare] Mount test_data/config.yaml as /data — see test-docker.ps1"
        exit 1
    fi
    log_info "Plain Docker mode: using pre-mounted /data/config.yaml"

    # Extract values needed for env propagation
    FTP_USERNAME=$(python3 -c "import yaml; c=yaml.safe_load(open('/data/config.yaml')); print(c.get('ftp',{}).get('username',''))" 2>/dev/null || echo "")
    FTP_PASSWORD=$(python3 -c "import yaml; c=yaml.safe_load(open('/data/config.yaml')); print(c.get('ftp',{}).get('password',''))" 2>/dev/null || echo "")
    TG_TOKEN=$(python3 -c "import yaml; c=yaml.safe_load(open('/data/config.yaml')); print(c.get('telegram',{}).get('bot_token',''))" 2>/dev/null || echo "")
fi

# ── Create directories ───────────────────────────────────────────────────────
log_info "Setting up directories..."
mkdir -p /share/scan_inbox/{scan_duplex,copy_duplex,scan_document,card_2in1,confirm,confirm_print,reject,test_print}
mkdir -p /share/scan_out
mkdir -p /data/checkpoints

# Copy checkpoint models to persistent storage on first run
if [ ! -f "/data/checkpoints/depth_anything_v2_vits_slim.onnx" ]; then
    if [ -d "/app/checkpoints" ] && [ "$(ls -A /app/checkpoints 2>/dev/null)" ]; then
        log_info "Copying checkpoint models to /data/checkpoints..."
        cp -r /app/checkpoints/* /data/checkpoints/
    fi
fi

# ── Propagate settings to S6 container environment for service scripts ───────
printf '%s' "/share/scan_inbox"  > /run/s6/container_environment/SCAN_INBOX_BASE
printf '%s' "/share/scan_inbox"  > /run/s6/container_environment/SCAN_INBOX_DIR
printf '%s' "/share/scan_out"    > /run/s6/container_environment/SCAN_OUTPUT_DIR
printf '%s' "${FTP_USERNAME:-}"  > /run/s6/container_environment/FTP_USERNAME
printf '%s' "${FTP_PASSWORD:-}"  > /run/s6/container_environment/FTP_PASSWORD
printf '%s' "${TG_TOKEN:-}"      > /run/s6/container_environment/SCAN_TELEGRAM_BOT_TOKEN
printf '%s' "1"                  > /run/s6/container_environment/PYTHONUNBUFFERED

log_info "Preparation complete. Starting services..."
