# ha-mic-filter/rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Real-time Microphone Filter
# No longer needed - Home Assistant provides PulseAudio via audio: true
# ==============================================================================

bashio::log.info "PulseAudio service disabled - using Home Assistant's audio system"

# Just sleep to keep the service alive but do nothing
# The virtual devices are created in the init script
exec sleep infinity