# ha-mic-filter/rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Real-time Microphone Filter
# Runs PulseAudio server
# ==============================================================================

bashio::log.info "Starting PulseAudio server..."

# Set PulseAudio configuration path
export PULSE_CONFIG_PATH=/tmp/pulse-config

# Start PulseAudio as system service with custom config
exec pulseaudio \
    --system \
    --disallow-exit \
    --daemonize=false \
    --log-target=stderr \
    --log-level=info \
    --use-pid-file=false \
    --disable-shm \
    --exit-idle-time=-1 \
    --file=/tmp/pulse-config/system.pa \
    --disable-memfd